# install/systemd/omnimind-mcp.service
[Unit]
Description=OmniMind MCP Servers
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/fahbrain/projects/omnimind
EnvironmentFile=/home/fahbrain/projects/omnimind/.env
ExecStart=/home/fahbrain/projects/omnimind/scripts/start_mcp_servers.sh
ExecStop=/usr/bin/pkill -f mcp
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
