# archive/legacy_scripts_20251216/omnimind-qdrant.service
[Unit]
Description=OmniMind Qdrant Vector Database
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/fahbrain/projects/omnimind
ExecStart=/usr/bin/docker-compose -f deploy/docker-compose.yml up qdrant
ExecStop=/usr/bin/docker-compose -f deploy/docker-compose.yml stop qdrant
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
