# config/systemd/mind.service
[Unit]
Description=OmniMind MCP System - Core Infrastructure
Documentation=https://omnimind.local/docs
After=network.target
Wants=mind-mcp-internal.service

[Service]
Type=simple
User=fahbrain
Group=fahbrain
WorkingDirectory=/home/fahbrain/projects/omnimind
Environment="PYTHONPATH=/home/fahbrain/projects/omnimind:/home/fahbrain/projects/omnimind/src"
Environment="OMNIMIND_ENV=production"
Environment="LOG_LEVEL=INFO"

# Core MCPs startup
ExecStart=/home/fahbrain/projects/omnimind/scripts/production/start_mcp_servers.sh

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=3

# Process management
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30

# Security
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=true
ReadWritePaths=/home/fahbrain/projects/omnimind/logs /home/fahbrain/projects/omnimind/data

# Resource limits
MemoryMax=4G
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=omnimind

[Install]
WantedBy=multi-user.target
