# scripts/production/deploy/omnimind-backup.service
[Unit]
Description=OmniMind Automated Backup Service
After=network.target
Requires=network.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/home/fahbrain/projects/omnimind/scripts/backup/automated_backup.sh
WorkingDirectory=/home/fahbrain/projects/omnimind
StandardOutput=journal
StandardError=journal
# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/backups /run/media/fahbrain/DEV_BRAIN_CLEAN
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
