# scripts/production/deploy/omnimind.service
[Unit]
Description=OmniMind Essential Services - Backend + Orchestrator (Fase 1: Essenciais)
Documentation=https://github.com/fabs-devbrain/OmniMind
After=network-online.target omnimind-qdrant.service
Wants=network-online.target
Before=omnimind-daemon.service omnimind-frontend.service omnimind-core.service

[Service]
Type=simple
User=fahbrain
Group=fahbrain
WorkingDirectory=/home/fahbrain/projects/omnimind

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="OMNIMIND_WORKSPACE=/home/fahbrain/projects/omnimind"
Environment="OMNIMIND_LOG_LEVEL=INFO"
Environment="PYTHONPATH=/home/fahbrain/projects/omnimind"
Environment="CUDA_HOME=/usr"
Environment="CUDA_VISIBLE_DEVICES=0"
Environment="CUDA_PATH=/usr"
Environment="LD_LIBRARY_PATH=/usr/lib/x86_64-linux-gnu"
Environment="PYTORCH_CUDA_ALLOC_CONF=backend:cudaMallocAsync"
Environment="OMNIMIND_MODE=production"

# Executar script oficial de inicialização
ExecStart=/bin/bash -c 'source /home/fahbrain/projects/omnimind/.venv/bin/activate && cd /home/fahbrain/projects/omnimind && ./scripts/canonical/system/start_omnimind_system.sh'

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitBurst=5
StartLimitInterval=300

# Resource limits
MemoryMax=4G
CPUQuota=200%
TasksMax=100

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/fahbrain/projects/omnimind/logs /home/fahbrain/projects/omnimind/data

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=omnimind

[Install]
WantedBy=multi-user.target
