# scripts/production/deploy/omnimind-backup.timer
[Unit]
Description=Run OmniMind backup daily
Requires=omnimind-backup.service

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
