# scripts/production/deploy/security-validation.timer
[Unit]
Description=Timer for OmniMind Security Validation

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
